# rootfs/etc/s6-overlay/s6-rc.d/omeka-classic-setup/up
/etc/s6-overlay/scripts/omeka-classic-setup.sh

